# deployment/production/supervisor/supervisor.d/laravel-worker.ini
[program:queue_default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work redis --queue=default --sleep=3 --tries=3 --timeout=120 --max-time=3600
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/%(program_name)s_%(process_num)02d.log

[program:queue_send_message]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work redis --queue=send_message --sleep=3 --tries=3 --timeout=120 --max-time=3600
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/%(program_name)s_%(process_num)02d.log

[program:queue_handle_register]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work redis --queue=handle_register --sleep=3 --tries=3 --timeout=120 --max-time=3600
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/%(program_name)s_%(process_num)02d.log

[program:scheduler]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan schedule:work
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/%(program_name)s_%(process_num)02d.log
